# com.openexchange.admin.contextrestore/packaging/debian/postrm
#! /bin/bash
# postrm script for open-xchange-admin
#
# see: dh_installdeb(1)

set -e

if [ -e /opt/open-xchange/etc/oxfunctions.sh ]; then
   . /opt/open-xchange/etc/oxfunctions.sh
else
   pname=$(basename $0 .postrm)
   echo "Warning: cannot correctly remove $pname because open-xchange-admin already gone."
   exit 0
fi


# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
       remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)


        ;;

       purge|remove)
       
       if ox_isrunning; then
	   invoke-rc.d open-xchange-admin restart
       fi

        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
